# zerotrust/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Execute Cloudflare Daemon
# ==============================================================================

TOKEN=$(bashio::config 'token')

bashio::log.info "Starting Cloudflare Services"
exec cloudflared tunnel run --force --token $TOKEN
bashio::log.info "Cloudflare Services Started"
