# zerotrust/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Execute Cloudflare Daemon
# ==============================================================================

TOKEN=$(bashio::config 'token')

bashio::log.info "Starting Cloudflare Services"
bashio::log.info "Fetching Data from Zero Trust Platform"
bashio::log.info "This addon is under maintenance by RACKSYNC"
exec cloudflared tunnel run --force --token $TOKEN

