# zerotrust/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Execute Cloudflare Daemon
# ==============================================================================

TOKEN=$(bashio::config 'token')

bashio::log.info "Starting Cloudflare Services"
bashio::log.info "Made with Love from Thailand, Bring to you by RACKSYNC 🇹🇭"

mkdir -p /root/.cloudflared/
cp -Rv /ssl/* /root/.cloudflared/

if bashio::config.true 'config'; then
    bashio::log.info "Fetching tunnel data from config.yaml 🔧"
    exec cloudflared tunnel --config /root/.cloudflared/config.yaml run
else
    bashio::log.info "Fetching tunnel data from Token! 🔑"
    exec cloudflared tunnel run --force --token $TOKEN
fi



