# zerotrust-pi/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio

TOKEN=$(bashio::config 'token')

bashio::log.info "Starting Cloudflare Services"
bashio::log.info "Made with Love from Thailand, Bring to you by RACKSYNC 🇹🇭"

mkdir -p /root/.cloudflared/


if bashio::config.true 'no_autoupdate'; then
   NOUPDATE="--no-autoupdate"
else
   NOUPDATE=" "
fi

FLAG=$(bashio::config 'flag')

if bashio::config.true 'config'; then
    bashio::log.info "Fetching tunnel data from config.yaml 🔧"
    exec cp -Rv /ssl/config.yaml /root/.cloudflared/
    exec cloudflared tunnel $FLAG $NOUPDATE --config /root/.cloudflared/config.yaml run
else
    bashio::log.info "Fetching tunnel data from Token! 🔑"
    exec cloudflared tunnel $FLAG $NOUPDATE run --force --token $TOKEN
fi



